# t_gen_derivatives.service
[Unit]
Description=autotrades Derivatives Generator
After=network.target

[Service]
Type=simple
User=autotrades
Group=www-data
UMask=0002
WorkingDirectory=/var/www/autotrades/scripts
ExecStart=/var/www/autotrades/venv/bin/python /var/www/autotrades/scripts/gen_derivatives.py
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
